# bash/deployment/config/movuex-1.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid root
chdir /opt/movuex/package

env NODE_ENV=production
env PORT=3000

exec npm start